# container-files/etc/supervisor.d/shellinabox.conf
[program:shellinabox]
directory=/tmp/
user=root
command=/usr/sbin/shellinaboxd -s /:LOGIN --disable-ssl
autostart=true
autorestart=true
redirect_stderr=true
